# api/api.service
[Unit]
Description=CPSC222 Python HTTP API
After=network.target
[Service]
ExecStart=/usr.bin/python3 /home/user/CPSC222/api/server.py
Type=exec
Restart=always
[Install]
WantedBy=default.target
